# .config/systemd/user/bgrnds.timer
[Unit]
Description=Periodically change backdrops

[Timer]
OnUnitActiveSec=15m
OnBootSec=15m

[Install]
WantedBy=timers.target
