# .config/systemd/user/bgrnds.service
[Unit]
Description=Sets random Gnome desktop background

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'gnome_desktop_rotate %h/media/bgrnd'
