# .config/systemd/user/bgrnds.service
[Unit]
Description=Sets random desktop & login backdrops

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'backdrop_rotate %h/media/bgrnd'
