# .config/systemd/user/logrotate.timer
[Unit]
Description=Regularly runs logrotate on jzacsh logs

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
