# .config/systemd/user/lxc-autostart.service
[Unit]
Description=Autostart Unpriveleged LXCs

[Service]
Type=simple
ExecStart=/bin/bash -c 'lxc-autostart'

[Install]
WantedBy=default.target
