# .config/systemd/user/upspin.service
[Unit]
Description=UpspinFS (Headless)

[Service]
ExecStart=/bin/bash -c "upspinfs %h/usr/share/upspin/"

[Install]
WantedBy=default.target
