# .config/systemd/user/urxvtd.service
[Unit]
Description=RXVT-Unicode Daemon

[Service]
ExecStart=/usr/bin/urxvtd --opendisplay

[Install]
WantedBy=default.target
