# .config/autostart/urxvtd.desktop
[Desktop Entry]
Comment=RXVT-Unicode terminal daemon
Name=urxvtd
Type=Application
Exec=/usr/bin/urxvtd --opendisplay
