# .config/autostart/xmomdap.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=Xmodmap
Comment=Run my xmodmap
Exec=/usr/bin/xmodmap .Xmodmap
StartupNotify=true
Terminal=false
Hidden=false
