# .host/src/host.wajinru/borgmatic.timer
[Unit]
Description=Scheduler of borgmatic backups

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
